# scripts/systemd/esg-pull-deploy.service
# Optional: run the deploy script manually on the server, e.g. after a one-off fix:
#   sudo systemctl start esg-pull-deploy.service
#
# Automatic deploys when GitHub changes should use GitHub Actions (push to main),
# which SSHs in and runs scripts/pull-and-deploy.sh — not a systemd timer.
#
# Install (once):
#   sudo cp scripts/systemd/esg-pull-deploy.service /etc/systemd/system/
#   sudo systemctl daemon-reload
# Do NOT enable this on a schedule; there is no timer unit.

[Unit]
Description=ESG — one-shot pull from GitHub and docker compose rebuild (manual or ad-hoc)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/esg
ExecStart=/root/esg/scripts/pull-and-deploy.sh

# No [Install] section — do not `enable` this unit. Use `systemctl start …` only when
# you need a one-off deploy; routine deploys use GitHub Actions on push.
